# oci-authenticator/flatpak-oci-authenticator.service.in
[Unit]
Description=flatpak oci authenticator

[Service]
BusName=org.flatpak.Authenticator.Oci
ExecStart=@libexecdir@/flatpak-oci-authenticator
Type=dbus
