# session-helper/flatpak-session-helper.service.in
[Unit]
Description=flatpak session helper

[Service]
BusName=org.freedesktop.Flatpak
ExecStart=@libexecdir@/flatpak-session-helper
Type=dbus
