# session-helper/xdg-app-session-helper.service.in
[Unit]
Description=xdg-app session helper

[Service]
BusName=org.freedesktop.XdgApp
ExecStart=@libexecdir@/xdg-app-session-helper
Type=dbus
