# sideload-repos-systemd/flatpak-sideload-usb-repo.path.in
# This unit is intended to be installed in the systemd user instance, and
# depends on flatpak-sideload-repos-dir.service being in the system instance
# and running first. The idea here is that we add any USB drive mounts to the
# appropriate directory so Flatpak can find and pull from them in case they
# have flatpaks on them, both when a new drive is inserted and at the start of
# the user session.
[Path]
PathExists=@media_dir@/%u
PathChanged=@media_dir@/%u

[Install]
WantedBy=default.target
