# document-portal/xdg-document-portal.service.in
[Unit]
Description=flatpak document portal service

[Service]
BusName=org.freedesktop.portal.Documents
ExecStart=@libexecdir@/xdg-document-portal
Type=dbus
