# portal/flatpak-portal.service.in
[Unit]
Description=flatpak portal

[Service]
BusName=org.freedesktop.portal.Flatpak
ExecStart=@libexecdir@/flatpak-portal
Type=dbus
