# system-helper/flatpak-system-helper.service.in
[Unit]
Description=flatpak system helper

[Service]
BusName=org.freedesktop.Flatpak.SystemHelper
Environment=XDG_DATA_DIRS=@localstatedir@/lib/flatpak/exports/share/:/usr/local/share/:/usr/share/
ExecStart=@libexecdir@/flatpak-system-helper
Type=dbus
IOSchedulingClass=idle
