# system-helper/flatpak-system-helper.service.in
[Unit]
Description=flatpak system helper

[Service]
BusName=org.freedesktop.Flatpak.SystemHelper
ExecStart=@libexecdir@/flatpak-system-helper
Type=dbus
