# system-helper/xdg-app-system-helper.service.in
[Unit]
Description=xdg-app system helper

[Service]
BusName=org.freedesktop.XdgApp.SystemHelper
ExecStart=@libexecdir@/xdg-app-system-helper
Type=dbus
